# node.conf
[program:node]
command=npm start
directory=/local/var/www/app
autostart=true
autorestart=true
